# resque-pool/templates/default/resque-pool.init.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          resque-pool
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts resque-pool for <%= @application %>
# Description:       starts resque-pool for <%= @application %>
### END INIT INFO

export SL_NAME=resque-pool_rails
NAME=resque-pool_<%= @application %>
DESC=resque-pool_<%= @application %>
PID=<%= @current_path %>/tmp/pids/resque-pool.pid

cd <%= @current_path %>

case "$1" in
  start)
    echo -n "Starting $DESC: "
    RAILS_ENV=<%= @rails_env %> bundle exec resque-pool --daemon --environment=<%= @rails_env %>
    echo "$NAME."
    ;;
  stop)
    echo -n "Stopping $DESC: "
    kill -QUIT `cat $PID`
    echo "$NAME."
    ;;
  restart)
    echo -n "Restarting $DESC: "
    kill -QUIT `cat $PID`
    sleep 1
    RAILS_ENV=<%= @rails_env %> bundle exec resque-pool --environment=<%= @rails_env %>
    echo "$NAME."
    ;;
  reload)
    echo -n "Reloading $DESC configuration: "
    kill -HUP `cat $PID`
    echo "$NAME."
    ;;
  *)
    echo "Usage: $NAME {start|stop|restart|reload}" >&2
    exit 1
    ;;
esac

exit 0
