# unicorn/templates/default/unicorn.init.erb
export SL_NAME=unicorn_rails
#! /bin/sh

### BEGIN INIT INFO
# Provides:          unicorn
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the unicorn web server for <%= @application %>
# Description:       starts unicorn for <%= @application %>
### END INIT INFO

<% options = "#{@app_path}/shared/config/unicorn.rb -E #{@rails_env} -D" %>
NAME=unicorn_<%= @application %>
DESC=unicorn_<%= @application %>
PID=<%= @app_path %>/shared/pids/unicorn_<%= @application %>.pid

cd <%= @current_path %>

case "$1" in
  start)
	echo -n "Starting $DESC: "
  bundle exec unicorn_rails -c <%= options %>
	echo "$NAME."
	;;
  stop)
	echo -n "Stopping $DESC: "
        kill -QUIT `cat $PID`
	echo "$NAME."
	;;
  restart)
	echo -n "Restarting $DESC: "
        kill -QUIT `cat $PID`
	sleep 1
  bundle exec unicorn_rails -c <%= options %>
	echo "$NAME."
	;;
  reload)
        echo -n "Reloading $DESC configuration: "
        kill -HUP `cat $PID`
        echo "$NAME."
        ;;
  *)
	echo "Usage: $NAME {start|stop|restart|reload}" >&2
	exit 1
	;;
esac

exit 0
